# {{cookiecutter.project_slug}}/deploy/supervisor/{{cookiecutter.project_slug}}_prd.ini
[program:tjapi]
command=/data/prd/{{cookiecutter.project_slug}}/env/bin/gunicorn -c gunicorn.conf.py -p gunicorn-{{cookiecutter.project_slug}}.pid  {{cookiecutter.project_slug}}config.wsgi
directory=/data/prd/{{cookiecutter.project_slug}}/src
user=pyer
autostart=true
autorestart=true
startretries=5
stdout_logfile=/data/prd/{{cookiecutter.project_slug}}/logs/supervisor_stdout.log
stderr_logfile=/data/prd/{{cookiecutter.project_slug}}/logs/supervisor_stderr.log
